# docker/supervisord.conf
[supervisord]
nodaemon=true

; Socket + interfaz RPC para que `supervisorctl status` funcione dentro del
; contenedor. Sin estas secciones supervisord corre los programas igual, pero
; el cliente falla con "unix:///var/run/supervisor.sock no such file".
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=php-fpm
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Corre las tareas programadas (backups spatie, avisos de vencimiento de
; licencia, backup:clean). Sin este proceso NADA de lo agendado se ejecuta en
; producción. schedule:work despierta cada minuto y lanza lo que toque.
[program:scheduler]
command=php /var/www/artisan schedule:work
user=www-data
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0